# atlas-agent/packaging/systemd/atlas-spatial-runtime.service
[Unit]
Description=Atlas Spatial Camera Runtime
After=docker.service
Requires=docker.service
ConditionPathExists=/etc/atlas-agent/spatial.env

[Service]
Type=simple
EnvironmentFile=/etc/atlas-agent/spatial.env
ExecStartPre=/usr/bin/docker image inspect --format={{.Id}} ${ATLAS_SPATIAL_CONTAINER_IMAGE}
ExecStartPre=-/usr/bin/docker rm -f ${ATLAS_SPATIAL_CONTAINER_NAME}
ExecStart=/usr/libexec/atlas-agent/atlas-spatial-container-run
ExecStop=/usr/bin/docker stop --timeout 10 ${ATLAS_SPATIAL_CONTAINER_NAME}
Restart=always
RestartSec=3
TimeoutStartSec=0
TimeoutStopSec=20
KillMode=process

[Install]
WantedBy=multi-user.target
